# config/ntpdate.cron
#!/bin/sh

/etc/init.d/ntp stop
/usr/sbin/ntpdate pool.ntp.org
/etc/init.d/ntp start
