# config/supervisor_node.conf
[program:node]
command=node unread_counts.js
directory=/home/sclay/newsblur/node
user=sclay
autostart=true
autorestart=true
#redirect_stderr=True
priority=991
stopsignal=HUP
stdout_logfile = /home/sclay/newsblur/logs/unread_counts.log
